# deploy/supervisord/sockjs.ini
[program:vcweb-sockjs]
directory=/opt/vcweb/
command=/opt/virtualenvs/vcweb/bin/python /opt/vcweb/vcweb/vcweb-sockjs.py 8882
user=nginx
autostart=true
redirect_stderr=true
