# deploy/supervisord/vcweb-uwsgi.ini
[program:vcweb-sockjs]
directory=/opt/vcweb/
command=/opt/virtualenvs/vcweb/bin/uwsgi --socket /opt/vcweb/uwsgi.sock --wsgi-file /opt/vcweb/vcweb/wsgi.py --master --processes 4 --threads 2 --stats 127.0.0.1:9191
user=apache
autostart=true
redirect_stderr=true
