# deploy/supervisord/vcweb-uwsgi.ini
[program:vcweb-uwsgi]
directory=/opt/vcweb/
command=/opt/virtualenvs/vcweb/bin/uwsgi /opt/vcweb/deploy/uwsgi/vcweb.ini
user=vcweb
autostart=true
redirect_stderr=true
