# deploy/supervisord.conf
[program:vcweb-sockjs]
directory=/opt/vcweb/
command=/opt/virtualenvs/vcweb/bin/python /opt/vcweb/vcweb/vcweb-sockjs.py 8882
autostart=true
redirect_stderr=true
