# deploy/supervisord.d/vcweb-uwsgi.ini
[program:vcweb-sockjs]
directory=/opt/vcweb/
command=/opt/virtualenvs/vcweb/bin/uwsgi --socket /opt/vcweb/vcweb.uwsgi.sock --wsgi-file /opt/vcweb/vcweb/wsgi.py
user=apache
autostart=true
redirect_stderr=true
